# os/linux/autostart/unclutter-autostart.desktop
[Desktop Entry]
Name=Unclutter
Exec=unclutter
Type=Application
NoDisplay=true
