# os/linux/autostart/xrandr-setup.desktop
[Desktop Entry]
Type=Application
Name=Fixmonitors
Exec=$HOME/.config/xrandr_setup.sh


